# deploy/topquaranta-web.service
[Unit]
Description=TopQuaranta - Public Web (Django)
After=network.target postgresql.service

[Service]
User=topquaranta
Group=topquaranta
WorkingDirectory=/home/topquaranta/app
EnvironmentFile=/home/topquaranta/app/.env
ExecStart=/home/topquaranta/app/.venv/bin/gunicorn topquaranta.wsgi:application --bind 127.0.0.1:8083 --workers 2 --timeout 120 --env DJANGO_SETTINGS_MODULE=topquaranta.settings.web_server
Restart=always

[Install]
WantedBy=multi-user.target
